# products.gallery.service
[Unit]
Description=products.gallery

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/home/pi/code/products.gallery/server/server
WorkingDirectory=/home/pi/code/products.gallery/server/

[Install]
WantedBy=multi-user.target