# products.gallery.service
[Unit]
Description=products.gallery

[Service]
Type=simple
TimeoutSec=900
ExecStartPre=/bin/sh -c 'cd ui rm -rf dist && rm -rf build && yarn build && mv build dist'
ExecStart=/home/pi/code/products.gallery/server/server
WorkingDirectory=/home/pi/code/products.gallery/

[Install]
WantedBy=multi-user.target